# .devcontainer/supervisord-novnc.conf
[supervisord]
nodaemon=true
logfile=/var/supervisor/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Xvfb: Virtual X server
[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_VNC_RESOLUTION)s x24
autostart=true
autorestart=true
stderr_logfile=/var/supervisor/xvfb.log
stdout_logfile=/var/supervisor/xvfb.log
priority=1

# Fluxbox: Lightweight window manager
[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=%(ENV_DISPLAY)s
autostart=true
autorestart=true
stderr_logfile=/var/supervisor/fluxbox.log
stdout_logfile=/var/supervisor/fluxbox.log
priority=2

# x11vnc: VNC server
[program:x11vnc]
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -passwd %(ENV_VNC_PASSWORD)s -forever -shared -xkb
autostart=true
autorestart=true
stderr_logfile=/var/supervisor/x11vnc.log
stdout_logfile=/var/supervisor/x11vnc.log
priority=3

# noVNC: Web-based VNC viewer
[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
stderr_logfile=/var/supervisor/novnc.log
stdout_logfile=/var/supervisor/novnc.log
priority=4
